# docker/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
silent=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:edu_system]
# supervisor启动命令
command = /root/mo/edu_system/venv/bin/gunicorn -c /root/mo/edu_system/gunicorn_config.p --workers 3 --bind unix:/run/gunicorn.sock edu_system.wsgi:application
directory = /root/mo/edu_system
autostart = true
autorestart = true
# 用于配置是否将子进程的标准错误（stderr）输出重定向到 stdout 的选项。当设置为 true 时，子进程的标准错误输出将被合并到标准输出中，这样可以方便地通过 stdout 日志来查看错误信息。
redirect_stderr = true
# 标准输出日志文件，容器中的日志文件映射到docker/logs/w5.log
stdout_logfile=/root/mo/edu_system/supervisor.log
# 日志文件最大容量
stdout_logfile_maxbytes = 20MB
# 日志文件最多备份数量
stdout_logfile_backups = 20